# services/motioneye/run
#!/usr/bin/env bash

set -euo pipefail
IFS=$'\n\t'

exec /sbin/setuser nobody /usr/local/bin/meyectl startserver -c /config/motioneye.conf
